# systemd/makepad-postgres-backup.timer
[Unit]
Description=Nightly Makepad PostgreSQL backup

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true

[Install]
WantedBy=timers.target
